# downloader-backend.service
[Unit]
Description=Downloader Backend Service
After=network.target

[Service]
# Change this to your user
User=root
# Change this to your working directory
WorkingDirectory=/opt/downloader/backend
# Path to the executable
# Run `swift build -c release --show-bin-path` to find the directory, then append "/Run"
ExecStart=/opt/downloader/backend/.build/release/Run serve --hostname 0.0.0.0 --port 8080
Restart=always
RestartSec=5
Environment="DATA_DIR=/opt/downloader/data"
# Add your API key here or in a .env file
# Environment="ONEFICHIER_API_KEY=your_key"

[Install]
WantedBy=multi-user.target
